# debian/postinst
#!/bin/bash
set -e

systemctl enable jpu
systemctl start jpu
