# usr/lib/systemd/system/jpu.service
###
 # Copyright 2022-2023 SPACEMIT. All rights reserved.
 # Use of this source code is governed by a BSD-style license
 # that can be found in the LICENSE file.
 # 
 # @Author: David(qiang.fu@spacemit.com)
 # @Date: 2023-09-20 13:41:02
 # @LastEditTime: 2023-09-20 14:05:04
 # @Description: 
###

[Unit]
Description=jpu service

[Service]
Type=oneshot
EnvironmentFile=/etc/init.d/jpu.sh
ExecStart=/etc/init.d/jpu.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
